# pi_scripts/scheduling/record-sound.service
Unit]
Description=Record 30 seconds of sound

[Service]
ExecStart=/path/to/record_sound.sh